# supervisord.conf
[supervisord]
nodaemon=true

[program:apache2]
command=/usr/sbin/apache2ctl -D FOREGROUND
autorestart=true

[program:ratchet]
command=php /var/www/html/servidorWebsocket.php
autorestart=true
